# scripts/fimwebui.service
[Unit]
Description=Fim web UI
After=network.target

[Service]
Type=simple
User=fim
WorkingDirectory=/opt/fim
Environment=FIM_DB_PATH=/opt/fim/data/fim.sqlite3
Environment=FIM_WEB_HOST=0.0.0.0
Environment=FIM_WEB_PORT=19992
ExecStart=/opt/fim/start_webui.sh
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
